# alpm/tests/db/local/audit-2.8.4-3/files
%FILES%
etc/
etc/audisp/
etc/audisp/audisp-remote.conf
etc/audisp/audispd.conf
etc/audisp/plugins.d/
etc/audisp/plugins.d/af_unix.conf
etc/audisp/plugins.d/au-remote.conf
etc/audisp/plugins.d/audispd-zos-remote.conf
etc/audisp/plugins.d/syslog.conf
etc/audisp/zos-remote.conf
etc/audit/
etc/audit/audit-stop.rules
etc/audit/auditd.conf
etc/libaudit.conf
usr/
usr/bin/
usr/bin/audisp-remote
usr/bin/audispd
usr/bin/audispd-zos-remote
usr/bin/auditctl
usr/bin/auditd
usr/bin/augenrules
usr/bin/aulast
usr/bin/aulastlog
usr/bin/aureport
usr/bin/ausearch
usr/bin/ausyscall
usr/bin/autrace
usr/bin/auvirt
usr/include/
usr/include/auparse-defs.h
usr/include/auparse.h
usr/include/libaudit.h
usr/lib/
usr/lib/libaudit.so
usr/lib/libaudit.so.1
usr/lib/libaudit.so.1.0.0
usr/lib/libauparse.so
usr/lib/libauparse.so.0
usr/lib/libauparse.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/audit.pc
usr/lib/pkgconfig/auparse.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/auditd.service
usr/share/
usr/share/aclocal/
usr/share/aclocal/audit.m4
usr/share/man/
usr/share/man/man3/
usr/share/man/man3/audit_add_rule_data.3.gz
usr/share/man/man3/audit_add_watch.3.gz
usr/share/man/man3/audit_delete_rule_data.3.gz
usr/share/man/man3/audit_detect_machine.3.gz
usr/share/man/man3/audit_encode_nv_string.3.gz
usr/share/man/man3/audit_get_reply.3.gz
usr/share/man/man3/audit_get_session.3.gz
usr/share/man/man3/audit_getloginuid.3.gz
usr/share/man/man3/audit_log_acct_message.3.gz
usr/share/man/man3/audit_log_semanage_message.3.gz
usr/share/man/man3/audit_log_user_avc_message.3.gz
usr/share/man/man3/audit_log_user_comm_message.3.gz
usr/share/man/man3/audit_log_user_command.3.gz
usr/share/man/man3/audit_log_user_message.3.gz
usr/share/man/man3/audit_open.3.gz
usr/share/man/man3/audit_request_rules_list_data.3.gz
usr/share/man/man3/audit_request_signal_info.3.gz
usr/share/man/man3/audit_request_status.3.gz
usr/share/man/man3/audit_set_backlog_limit.3.gz
usr/share/man/man3/audit_set_backlog_wait_time.3.gz
usr/share/man/man3/audit_set_enabled.3.gz
usr/share/man/man3/audit_set_failure.3.gz
usr/share/man/man3/audit_set_pid.3.gz
usr/share/man/man3/audit_set_rate_limit.3.gz
usr/share/man/man3/audit_setloginuid.3.gz
usr/share/man/man3/audit_update_watch_perms.3.gz
usr/share/man/man3/auparse_add_callback.3.gz
usr/share/man/man3/auparse_destroy.3.gz
usr/share/man/man3/auparse_feed.3.gz
usr/share/man/man3/auparse_feed_age_events.3.gz
usr/share/man/man3/auparse_feed_has_data.3.gz
usr/share/man/man3/auparse_find_field.3.gz
usr/share/man/man3/auparse_find_field_next.3.gz
usr/share/man/man3/auparse_first_field.3.gz
usr/share/man/man3/auparse_first_record.3.gz
usr/share/man/man3/auparse_flush_feed.3.gz
usr/share/man/man3/auparse_get_field_int.3.gz
usr/share/man/man3/auparse_get_field_name.3.gz
usr/share/man/man3/auparse_get_field_num.3.gz
usr/share/man/man3/auparse_get_field_str.3.gz
usr/share/man/man3/auparse_get_field_type.3.gz
usr/share/man/man3/auparse_get_filename.3.gz
usr/share/man/man3/auparse_get_line_number.3.gz
usr/share/man/man3/auparse_get_milli.3.gz
usr/share/man/man3/auparse_get_node.3.gz
usr/share/man/man3/auparse_get_num_fields.3.gz
usr/share/man/man3/auparse_get_num_records.3.gz
usr/share/man/man3/auparse_get_record_num.3.gz
usr/share/man/man3/auparse_get_record_text.3.gz
usr/share/man/man3/auparse_get_serial.3.gz
usr/share/man/man3/auparse_get_time.3.gz
usr/share/man/man3/auparse_get_timestamp.3.gz
usr/share/man/man3/auparse_get_type.3.gz
usr/share/man/man3/auparse_get_type_name.3.gz
usr/share/man/man3/auparse_goto_field_num.3.gz
usr/share/man/man3/auparse_goto_record_num.3.gz
usr/share/man/man3/auparse_init.3.gz
usr/share/man/man3/auparse_interpret_field.3.gz
usr/share/man/man3/auparse_next_event.3.gz
usr/share/man/man3/auparse_next_field.3.gz
usr/share/man/man3/auparse_next_record.3.gz
usr/share/man/man3/auparse_node_compare.3.gz
usr/share/man/man3/auparse_normalize.3.gz
usr/share/man/man3/auparse_normalize_functions.3.gz
usr/share/man/man3/auparse_reset.3.gz
usr/share/man/man3/auparse_set_escape_mode.3.gz
usr/share/man/man3/auparse_timestamp_compare.3.gz
usr/share/man/man3/ausearch_add_expression.3.gz
usr/share/man/man3/ausearch_add_interpreted_item.3.gz
usr/share/man/man3/ausearch_add_item.3.gz
usr/share/man/man3/ausearch_add_regex.3.gz
usr/share/man/man3/ausearch_add_timestamp_item.3.gz
usr/share/man/man3/ausearch_add_timestamp_item_ex.3.gz
usr/share/man/man3/ausearch_clear.3.gz
usr/share/man/man3/ausearch_next_event.3.gz
usr/share/man/man3/ausearch_set_stop.3.gz
usr/share/man/man3/get_auditfail_action.3.gz
usr/share/man/man3/set_aumessage_mode.3.gz
usr/share/man/man5/
usr/share/man/man5/audisp-remote.conf.5.gz
usr/share/man/man5/audispd.conf.5.gz
usr/share/man/man5/auditd.conf.5.gz
usr/share/man/man5/ausearch-expression.5.gz
usr/share/man/man5/libaudit.conf.5.gz
usr/share/man/man5/zos-remote.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/audit.rules.7.gz
usr/share/man/man8/
usr/share/man/man8/audisp-remote.8.gz
usr/share/man/man8/audispd-zos-remote.8.gz
usr/share/man/man8/audispd.8.gz
usr/share/man/man8/auditctl.8.gz
usr/share/man/man8/auditd.8.gz
usr/share/man/man8/augenrules.8.gz
usr/share/man/man8/aulast.8.gz
usr/share/man/man8/aulastlog.8.gz
usr/share/man/man8/aureport.8.gz
usr/share/man/man8/ausearch.8.gz
usr/share/man/man8/ausyscall.8.gz
usr/share/man/man8/autrace.8.gz
usr/share/man/man8/auvirt.8.gz
var/
var/log/
var/log/audit/

%BACKUP%
etc/libaudit.conf	(null)
etc/audit/audit-stop.rules	(null)
etc/audit/auditd.conf	(null)
etc/audisp/audispd.conf	(null)
etc/audisp/audisp-remote.conf	(null)
etc/audisp/zos-remote.conf	(null)
etc/audisp/plugins.d/af_unix.conf	(null)
etc/audisp/plugins.d/audispd-zos-remote.conf	(null)
etc/audisp/plugins.d/au-remote.conf	(null)
etc/audisp/plugins.d/syslog.conf	(null)

