# alpm/tests/db/local/shadow-4.6-2/files
%FILES%
etc/
etc/default/
etc/default/useradd
etc/login.defs
etc/pam.d/
etc/pam.d/chage
etc/pam.d/chgpasswd
etc/pam.d/chpasswd
etc/pam.d/groupadd
etc/pam.d/groupdel
etc/pam.d/groupmems
etc/pam.d/groupmod
etc/pam.d/newusers
etc/pam.d/passwd
etc/pam.d/shadow
etc/pam.d/useradd
etc/pam.d/userdel
etc/pam.d/usermod
usr/
usr/bin/
usr/bin/chage
usr/bin/chgpasswd
usr/bin/chpasswd
usr/bin/expiry
usr/bin/faillog
usr/bin/gpasswd
usr/bin/groupadd
usr/bin/groupdel
usr/bin/groupmems
usr/bin/groupmod
usr/bin/groups
usr/bin/grpck
usr/bin/grpconv
usr/bin/grpunconv
usr/bin/lastlog
usr/bin/newgidmap
usr/bin/newuidmap
usr/bin/newusers
usr/bin/passwd
usr/bin/pwck
usr/bin/pwconv
usr/bin/pwunconv
usr/bin/sg
usr/bin/useradd
usr/bin/userdel
usr/bin/usermod
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shadow.service
usr/lib/systemd/system/shadow.timer
usr/lib/systemd/system/timers.target.wants/
usr/lib/systemd/system/timers.target.wants/shadow.timer
usr/share/
usr/share/licenses/
usr/share/licenses/shadow/
usr/share/licenses/shadow/LICENSE
usr/share/locale/
usr/share/locale/bs/
usr/share/locale/bs/LC_MESSAGES/
usr/share/locale/bs/LC_MESSAGES/shadow.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/shadow.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/shadow.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/shadow.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/shadow.mo
usr/share/locale/dz/
usr/share/locale/dz/LC_MESSAGES/
usr/share/locale/dz/LC_MESSAGES/shadow.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/shadow.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/shadow.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/shadow.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/shadow.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/shadow.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/shadow.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/shadow.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/shadow.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/shadow.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/shadow.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/shadow.mo
usr/share/locale/kk/
usr/share/locale/kk/LC_MESSAGES/
usr/share/locale/kk/LC_MESSAGES/shadow.mo
usr/share/locale/km/
usr/share/locale/km/LC_MESSAGES/
usr/share/locale/km/LC_MESSAGES/shadow.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/shadow.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/shadow.mo
usr/share/locale/ne/
usr/share/locale/ne/LC_MESSAGES/
usr/share/locale/ne/LC_MESSAGES/shadow.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/shadow.mo
usr/share/locale/nn/
usr/share/locale/nn/LC_MESSAGES/
usr/share/locale/nn/LC_MESSAGES/shadow.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/shadow.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/shadow.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/shadow.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/shadow.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/shadow.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/shadow.mo
usr/share/locale/sq/
usr/share/locale/sq/LC_MESSAGES/
usr/share/locale/sq/LC_MESSAGES/shadow.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/shadow.mo
usr/share/locale/tl/
usr/share/locale/tl/LC_MESSAGES/
usr/share/locale/tl/LC_MESSAGES/shadow.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/shadow.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/shadow.mo
usr/share/locale/vi/
usr/share/locale/vi/LC_MESSAGES/
usr/share/locale/vi/LC_MESSAGES/shadow.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/shadow.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/shadow.mo
usr/share/man/
usr/share/man/cs/
usr/share/man/cs/man1/
usr/share/man/cs/man1/expiry.1.gz
usr/share/man/cs/man1/gpasswd.1.gz
usr/share/man/cs/man1/groups.1.gz
usr/share/man/cs/man5/
usr/share/man/cs/man5/faillog.5.gz
usr/share/man/cs/man5/gshadow.5.gz
usr/share/man/cs/man5/passwd.5.gz
usr/share/man/cs/man5/shadow.5.gz
usr/share/man/cs/man8/
usr/share/man/cs/man8/faillog.8.gz
usr/share/man/cs/man8/groupadd.8.gz
usr/share/man/cs/man8/groupdel.8.gz
usr/share/man/cs/man8/groupmod.8.gz
usr/share/man/cs/man8/grpck.8.gz
usr/share/man/cs/man8/lastlog.8.gz
usr/share/man/da/
usr/share/man/da/man1/
usr/share/man/da/man1/groups.1.gz
usr/share/man/da/man1/sg.1.gz
usr/share/man/da/man5/
usr/share/man/da/man5/gshadow.5.gz
usr/share/man/da/man8/
usr/share/man/da/man8/groupdel.8.gz
usr/share/man/de/
usr/share/man/de/man1/
usr/share/man/de/man1/chage.1.gz
usr/share/man/de/man1/expiry.1.gz
usr/share/man/de/man1/gpasswd.1.gz
usr/share/man/de/man1/groups.1.gz
usr/share/man/de/man1/passwd.1.gz
usr/share/man/de/man1/sg.1.gz
usr/share/man/de/man3/
usr/share/man/de/man3/getspnam.3.gz
usr/share/man/de/man3/shadow.3.gz
usr/share/man/de/man5/
usr/share/man/de/man5/faillog.5.gz
usr/share/man/de/man5/gshadow.5.gz
usr/share/man/de/man5/login.defs.5.gz
usr/share/man/de/man5/passwd.5.gz
usr/share/man/de/man5/shadow.5.gz
usr/share/man/de/man5/suauth.5.gz
usr/share/man/de/man8/
usr/share/man/de/man8/chgpasswd.8.gz
usr/share/man/de/man8/chpasswd.8.gz
usr/share/man/de/man8/faillog.8.gz
usr/share/man/de/man8/groupadd.8.gz
usr/share/man/de/man8/groupdel.8.gz
usr/share/man/de/man8/groupmems.8.gz
usr/share/man/de/man8/groupmod.8.gz
usr/share/man/de/man8/grpck.8.gz
usr/share/man/de/man8/grpconv.8.gz
usr/share/man/de/man8/grpunconv.8.gz
usr/share/man/de/man8/lastlog.8.gz
usr/share/man/de/man8/newusers.8.gz
usr/share/man/de/man8/pwck.8.gz
usr/share/man/de/man8/pwconv.8.gz
usr/share/man/de/man8/pwunconv.8.gz
usr/share/man/de/man8/useradd.8.gz
usr/share/man/de/man8/userdel.8.gz
usr/share/man/de/man8/usermod.8.gz
usr/share/man/fr/
usr/share/man/fr/man1/
usr/share/man/fr/man1/chage.1.gz
usr/share/man/fr/man1/expiry.1.gz
usr/share/man/fr/man1/gpasswd.1.gz
usr/share/man/fr/man1/groups.1.gz
usr/share/man/fr/man1/newgidmap.1.gz
usr/share/man/fr/man1/newuidmap.1.gz
usr/share/man/fr/man1/passwd.1.gz
usr/share/man/fr/man1/sg.1.gz
usr/share/man/fr/man3/
usr/share/man/fr/man3/getspnam.3.gz
usr/share/man/fr/man3/shadow.3.gz
usr/share/man/fr/man5/
usr/share/man/fr/man5/faillog.5.gz
usr/share/man/fr/man5/gshadow.5.gz
usr/share/man/fr/man5/login.defs.5.gz
usr/share/man/fr/man5/passwd.5.gz
usr/share/man/fr/man5/shadow.5.gz
usr/share/man/fr/man5/suauth.5.gz
usr/share/man/fr/man5/subgid.5.gz
usr/share/man/fr/man5/subuid.5.gz
usr/share/man/fr/man8/
usr/share/man/fr/man8/chgpasswd.8.gz
usr/share/man/fr/man8/chpasswd.8.gz
usr/share/man/fr/man8/faillog.8.gz
usr/share/man/fr/man8/groupadd.8.gz
usr/share/man/fr/man8/groupdel.8.gz
usr/share/man/fr/man8/groupmems.8.gz
usr/share/man/fr/man8/groupmod.8.gz
usr/share/man/fr/man8/grpck.8.gz
usr/share/man/fr/man8/grpconv.8.gz
usr/share/man/fr/man8/grpunconv.8.gz
usr/share/man/fr/man8/lastlog.8.gz
usr/share/man/fr/man8/newusers.8.gz
usr/share/man/fr/man8/pwck.8.gz
usr/share/man/fr/man8/pwconv.8.gz
usr/share/man/fr/man8/pwunconv.8.gz
usr/share/man/fr/man8/useradd.8.gz
usr/share/man/fr/man8/userdel.8.gz
usr/share/man/fr/man8/usermod.8.gz
usr/share/man/hu/
usr/share/man/hu/man1/
usr/share/man/hu/man1/gpasswd.1.gz
usr/share/man/hu/man1/groups.1.gz
usr/share/man/hu/man1/passwd.1.gz
usr/share/man/hu/man1/sg.1.gz
usr/share/man/hu/man5/
usr/share/man/hu/man5/passwd.5.gz
usr/share/man/hu/man8/
usr/share/man/hu/man8/lastlog.8.gz
usr/share/man/id/
usr/share/man/id/man8/
usr/share/man/id/man8/useradd.8.gz
usr/share/man/it/
usr/share/man/it/man1/
usr/share/man/it/man1/chage.1.gz
usr/share/man/it/man1/expiry.1.gz
usr/share/man/it/man1/gpasswd.1.gz
usr/share/man/it/man1/groups.1.gz
usr/share/man/it/man1/passwd.1.gz
usr/share/man/it/man1/sg.1.gz
usr/share/man/it/man3/
usr/share/man/it/man3/getspnam.3.gz
usr/share/man/it/man3/shadow.3.gz
usr/share/man/it/man5/
usr/share/man/it/man5/faillog.5.gz
usr/share/man/it/man5/gshadow.5.gz
usr/share/man/it/man5/login.defs.5.gz
usr/share/man/it/man5/passwd.5.gz
usr/share/man/it/man5/shadow.5.gz
usr/share/man/it/man5/suauth.5.gz
usr/share/man/it/man8/
usr/share/man/it/man8/chgpasswd.8.gz
usr/share/man/it/man8/chpasswd.8.gz
usr/share/man/it/man8/faillog.8.gz
usr/share/man/it/man8/groupadd.8.gz
usr/share/man/it/man8/groupdel.8.gz
usr/share/man/it/man8/groupmems.8.gz
usr/share/man/it/man8/groupmod.8.gz
usr/share/man/it/man8/grpck.8.gz
usr/share/man/it/man8/grpconv.8.gz
usr/share/man/it/man8/grpunconv.8.gz
usr/share/man/it/man8/lastlog.8.gz
usr/share/man/it/man8/newusers.8.gz
usr/share/man/it/man8/pwck.8.gz
usr/share/man/it/man8/pwconv.8.gz
usr/share/man/it/man8/pwunconv.8.gz
usr/share/man/it/man8/useradd.8.gz
usr/share/man/it/man8/userdel.8.gz
usr/share/man/it/man8/usermod.8.gz
usr/share/man/ja/
usr/share/man/ja/man1/
usr/share/man/ja/man1/chage.1.gz
usr/share/man/ja/man1/expiry.1.gz
usr/share/man/ja/man1/gpasswd.1.gz
usr/share/man/ja/man1/groups.1.gz
usr/share/man/ja/man1/passwd.1.gz
usr/share/man/ja/man1/sg.1.gz
usr/share/man/ja/man5/
usr/share/man/ja/man5/faillog.5.gz
usr/share/man/ja/man5/login.defs.5.gz
usr/share/man/ja/man5/passwd.5.gz
usr/share/man/ja/man5/shadow.5.gz
usr/share/man/ja/man5/suauth.5.gz
usr/share/man/ja/man8/
usr/share/man/ja/man8/chpasswd.8.gz
usr/share/man/ja/man8/faillog.8.gz
usr/share/man/ja/man8/groupadd.8.gz
usr/share/man/ja/man8/groupdel.8.gz
usr/share/man/ja/man8/groupmod.8.gz
usr/share/man/ja/man8/grpck.8.gz
usr/share/man/ja/man8/grpconv.8.gz
usr/share/man/ja/man8/grpunconv.8.gz
usr/share/man/ja/man8/lastlog.8.gz
usr/share/man/ja/man8/newusers.8.gz
usr/share/man/ja/man8/pwck.8.gz
usr/share/man/ja/man8/pwconv.8.gz
usr/share/man/ja/man8/pwunconv.8.gz
usr/share/man/ja/man8/useradd.8.gz
usr/share/man/ja/man8/userdel.8.gz
usr/share/man/ja/man8/usermod.8.gz
usr/share/man/ko/
usr/share/man/ko/man1/
usr/share/man/ko/man1/groups.1.gz
usr/share/man/ko/man5/
usr/share/man/ko/man5/passwd.5.gz
usr/share/man/man1/
usr/share/man/man1/chage.1.gz
usr/share/man/man1/expiry.1.gz
usr/share/man/man1/gpasswd.1.gz
usr/share/man/man1/groups.1.gz
usr/share/man/man1/newgidmap.1.gz
usr/share/man/man1/newuidmap.1.gz
usr/share/man/man1/passwd.1.gz
usr/share/man/man1/sg.1.gz
usr/share/man/man3/
usr/share/man/man3/getspnam.3.gz
usr/share/man/man3/shadow.3.gz
usr/share/man/man5/
usr/share/man/man5/faillog.5.gz
usr/share/man/man5/gshadow.5.gz
usr/share/man/man5/login.defs.5.gz
usr/share/man/man5/passwd.5.gz
usr/share/man/man5/shadow.5.gz
usr/share/man/man5/suauth.5.gz
usr/share/man/man5/subgid.5.gz
usr/share/man/man5/subuid.5.gz
usr/share/man/man8/
usr/share/man/man8/chgpasswd.8.gz
usr/share/man/man8/chpasswd.8.gz
usr/share/man/man8/faillog.8.gz
usr/share/man/man8/groupadd.8.gz
usr/share/man/man8/groupdel.8.gz
usr/share/man/man8/groupmems.8.gz
usr/share/man/man8/groupmod.8.gz
usr/share/man/man8/grpck.8.gz
usr/share/man/man8/grpconv.8.gz
usr/share/man/man8/grpunconv.8.gz
usr/share/man/man8/lastlog.8.gz
usr/share/man/man8/newusers.8.gz
usr/share/man/man8/pwck.8.gz
usr/share/man/man8/pwconv.8.gz
usr/share/man/man8/pwunconv.8.gz
usr/share/man/man8/useradd.8.gz
usr/share/man/man8/userdel.8.gz
usr/share/man/man8/usermod.8.gz
usr/share/man/pl/
usr/share/man/pl/man1/
usr/share/man/pl/man1/chage.1.gz
usr/share/man/pl/man1/expiry.1.gz
usr/share/man/pl/man1/groups.1.gz
usr/share/man/pl/man1/sg.1.gz
usr/share/man/pl/man3/
usr/share/man/pl/man3/getspnam.3.gz
usr/share/man/pl/man3/shadow.3.gz
usr/share/man/pl/man5/
usr/share/man/pl/man5/faillog.5.gz
usr/share/man/pl/man8/
usr/share/man/pl/man8/faillog.8.gz
usr/share/man/pl/man8/groupadd.8.gz
usr/share/man/pl/man8/groupdel.8.gz
usr/share/man/pl/man8/groupmems.8.gz
usr/share/man/pl/man8/groupmod.8.gz
usr/share/man/pl/man8/grpck.8.gz
usr/share/man/pl/man8/lastlog.8.gz
usr/share/man/pl/man8/userdel.8.gz
usr/share/man/pl/man8/usermod.8.gz
usr/share/man/pt_BR/
usr/share/man/pt_BR/man1/
usr/share/man/pt_BR/man1/gpasswd.1.gz
usr/share/man/pt_BR/man5/
usr/share/man/pt_BR/man5/passwd.5.gz
usr/share/man/pt_BR/man5/shadow.5.gz
usr/share/man/pt_BR/man8/
usr/share/man/pt_BR/man8/groupadd.8.gz
usr/share/man/pt_BR/man8/groupdel.8.gz
usr/share/man/pt_BR/man8/groupmod.8.gz
usr/share/man/ru/
usr/share/man/ru/man1/
usr/share/man/ru/man1/chage.1.gz
usr/share/man/ru/man1/expiry.1.gz
usr/share/man/ru/man1/gpasswd.1.gz
usr/share/man/ru/man1/groups.1.gz
usr/share/man/ru/man1/passwd.1.gz
usr/share/man/ru/man1/sg.1.gz
usr/share/man/ru/man3/
usr/share/man/ru/man3/getspnam.3.gz
usr/share/man/ru/man3/shadow.3.gz
usr/share/man/ru/man5/
usr/share/man/ru/man5/faillog.5.gz
usr/share/man/ru/man5/gshadow.5.gz
usr/share/man/ru/man5/login.defs.5.gz
usr/share/man/ru/man5/passwd.5.gz
usr/share/man/ru/man5/shadow.5.gz
usr/share/man/ru/man5/suauth.5.gz
usr/share/man/ru/man8/
usr/share/man/ru/man8/chgpasswd.8.gz
usr/share/man/ru/man8/chpasswd.8.gz
usr/share/man/ru/man8/faillog.8.gz
usr/share/man/ru/man8/groupadd.8.gz
usr/share/man/ru/man8/groupdel.8.gz
usr/share/man/ru/man8/groupmems.8.gz
usr/share/man/ru/man8/groupmod.8.gz
usr/share/man/ru/man8/grpck.8.gz
usr/share/man/ru/man8/grpconv.8.gz
usr/share/man/ru/man8/grpunconv.8.gz
usr/share/man/ru/man8/lastlog.8.gz
usr/share/man/ru/man8/newusers.8.gz
usr/share/man/ru/man8/pwck.8.gz
usr/share/man/ru/man8/pwconv.8.gz
usr/share/man/ru/man8/pwunconv.8.gz
usr/share/man/ru/man8/useradd.8.gz
usr/share/man/ru/man8/userdel.8.gz
usr/share/man/ru/man8/usermod.8.gz
usr/share/man/sv/
usr/share/man/sv/man1/
usr/share/man/sv/man1/chage.1.gz
usr/share/man/sv/man1/expiry.1.gz
usr/share/man/sv/man1/groups.1.gz
usr/share/man/sv/man1/passwd.1.gz
usr/share/man/sv/man1/sg.1.gz
usr/share/man/sv/man3/
usr/share/man/sv/man3/getspnam.3.gz
usr/share/man/sv/man3/shadow.3.gz
usr/share/man/sv/man5/
usr/share/man/sv/man5/faillog.5.gz
usr/share/man/sv/man5/gshadow.5.gz
usr/share/man/sv/man5/passwd.5.gz
usr/share/man/sv/man5/suauth.5.gz
usr/share/man/sv/man8/
usr/share/man/sv/man8/faillog.8.gz
usr/share/man/sv/man8/groupadd.8.gz
usr/share/man/sv/man8/groupdel.8.gz
usr/share/man/sv/man8/groupmems.8.gz
usr/share/man/sv/man8/groupmod.8.gz
usr/share/man/sv/man8/grpck.8.gz
usr/share/man/sv/man8/lastlog.8.gz
usr/share/man/sv/man8/pwck.8.gz
usr/share/man/sv/man8/userdel.8.gz
usr/share/man/tr/
usr/share/man/tr/man1/
usr/share/man/tr/man1/chage.1.gz
usr/share/man/tr/man1/passwd.1.gz
usr/share/man/tr/man5/
usr/share/man/tr/man5/passwd.5.gz
usr/share/man/tr/man5/shadow.5.gz
usr/share/man/tr/man8/
usr/share/man/tr/man8/groupadd.8.gz
usr/share/man/tr/man8/groupdel.8.gz
usr/share/man/tr/man8/groupmod.8.gz
usr/share/man/tr/man8/useradd.8.gz
usr/share/man/tr/man8/userdel.8.gz
usr/share/man/tr/man8/usermod.8.gz
usr/share/man/zh_CN/
usr/share/man/zh_CN/man1/
usr/share/man/zh_CN/man1/chage.1.gz
usr/share/man/zh_CN/man1/expiry.1.gz
usr/share/man/zh_CN/man1/gpasswd.1.gz
usr/share/man/zh_CN/man1/groups.1.gz
usr/share/man/zh_CN/man1/passwd.1.gz
usr/share/man/zh_CN/man1/sg.1.gz
usr/share/man/zh_CN/man3/
usr/share/man/zh_CN/man3/getspnam.3.gz
usr/share/man/zh_CN/man3/shadow.3.gz
usr/share/man/zh_CN/man5/
usr/share/man/zh_CN/man5/faillog.5.gz
usr/share/man/zh_CN/man5/gshadow.5.gz
usr/share/man/zh_CN/man5/login.defs.5.gz
usr/share/man/zh_CN/man5/passwd.5.gz
usr/share/man/zh_CN/man5/shadow.5.gz
usr/share/man/zh_CN/man5/suauth.5.gz
usr/share/man/zh_CN/man8/
usr/share/man/zh_CN/man8/chgpasswd.8.gz
usr/share/man/zh_CN/man8/chpasswd.8.gz
usr/share/man/zh_CN/man8/faillog.8.gz
usr/share/man/zh_CN/man8/groupadd.8.gz
usr/share/man/zh_CN/man8/groupdel.8.gz
usr/share/man/zh_CN/man8/groupmems.8.gz
usr/share/man/zh_CN/man8/groupmod.8.gz
usr/share/man/zh_CN/man8/grpck.8.gz
usr/share/man/zh_CN/man8/grpconv.8.gz
usr/share/man/zh_CN/man8/grpunconv.8.gz
usr/share/man/zh_CN/man8/lastlog.8.gz
usr/share/man/zh_CN/man8/newusers.8.gz
usr/share/man/zh_CN/man8/pwck.8.gz
usr/share/man/zh_CN/man8/pwconv.8.gz
usr/share/man/zh_CN/man8/pwunconv.8.gz
usr/share/man/zh_CN/man8/useradd.8.gz
usr/share/man/zh_CN/man8/userdel.8.gz
usr/share/man/zh_CN/man8/usermod.8.gz
usr/share/man/zh_TW/
usr/share/man/zh_TW/man5/
usr/share/man/zh_TW/man5/passwd.5.gz
usr/share/man/zh_TW/man8/
usr/share/man/zh_TW/man8/chpasswd.8.gz
usr/share/man/zh_TW/man8/groupadd.8.gz
usr/share/man/zh_TW/man8/groupdel.8.gz
usr/share/man/zh_TW/man8/groupmod.8.gz
usr/share/man/zh_TW/man8/useradd.8.gz
usr/share/man/zh_TW/man8/userdel.8.gz
usr/share/man/zh_TW/man8/usermod.8.gz

%BACKUP%
etc/login.defs	(null)
etc/pam.d/chage	(null)
etc/pam.d/passwd	(null)
etc/pam.d/shadow	(null)
etc/pam.d/useradd	(null)
etc/pam.d/usermod	(null)
etc/pam.d/userdel	(null)
etc/pam.d/chpasswd	(null)
etc/pam.d/newusers	(null)
etc/pam.d/groupadd	(null)
etc/pam.d/groupdel	(null)
etc/pam.d/groupmod	(null)
etc/pam.d/chgpasswd	(null)
etc/pam.d/groupmems	(null)
etc/default/useradd	(null)

