# alpm/tests/db/local/dbus-1.12.12-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/session.conf
etc/dbus-1/system.conf
usr/
usr/bin/
usr/bin/dbus-cleanup-sockets
usr/bin/dbus-daemon
usr/bin/dbus-launch
usr/bin/dbus-monitor
usr/bin/dbus-run-session
usr/bin/dbus-send
usr/bin/dbus-test-tool
usr/bin/dbus-update-activation-environment
usr/bin/dbus-uuidgen
usr/include/
usr/include/dbus-1.0/
usr/include/dbus-1.0/dbus/
usr/include/dbus-1.0/dbus/dbus-address.h
usr/include/dbus-1.0/dbus/dbus-bus.h
usr/include/dbus-1.0/dbus/dbus-connection.h
usr/include/dbus-1.0/dbus/dbus-errors.h
usr/include/dbus-1.0/dbus/dbus-macros.h
usr/include/dbus-1.0/dbus/dbus-memory.h
usr/include/dbus-1.0/dbus/dbus-message.h
usr/include/dbus-1.0/dbus/dbus-misc.h
usr/include/dbus-1.0/dbus/dbus-pending-call.h
usr/include/dbus-1.0/dbus/dbus-protocol.h
usr/include/dbus-1.0/dbus/dbus-server.h
usr/include/dbus-1.0/dbus/dbus-shared.h
usr/include/dbus-1.0/dbus/dbus-signature.h
usr/include/dbus-1.0/dbus/dbus-syntax.h
usr/include/dbus-1.0/dbus/dbus-threads.h
usr/include/dbus-1.0/dbus/dbus-types.h
usr/include/dbus-1.0/dbus/dbus.h
usr/lib/
usr/lib/cmake/
usr/lib/cmake/DBus1/
usr/lib/cmake/DBus1/DBus1Config.cmake
usr/lib/cmake/DBus1/DBus1ConfigVersion.cmake
usr/lib/dbus-1.0/
usr/lib/dbus-1.0/dbus-daemon-launch-helper
usr/lib/dbus-1.0/include/
usr/lib/dbus-1.0/include/dbus/
usr/lib/dbus-1.0/include/dbus/dbus-arch-deps.h
usr/lib/libdbus-1.so
usr/lib/libdbus-1.so.3
usr/lib/libdbus-1.so.3.19.9
usr/lib/pkgconfig/
usr/lib/pkgconfig/dbus-1.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dbus.service
usr/lib/systemd/system/dbus.socket
usr/lib/systemd/system/multi-user.target.wants/
usr/lib/systemd/system/multi-user.target.wants/dbus.service
usr/lib/systemd/system/sockets.target.wants/
usr/lib/systemd/system/sockets.target.wants/dbus.socket
usr/lib/systemd/user/
usr/lib/systemd/user/dbus.service
usr/lib/systemd/user/dbus.socket
usr/lib/systemd/user/sockets.target.wants/
usr/lib/systemd/user/sockets.target.wants/dbus.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/dbus.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/dbus.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/session.conf
usr/share/dbus-1/session.d/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system.conf
usr/share/dbus-1/system.d/
usr/share/licenses/
usr/share/licenses/dbus/
usr/share/licenses/dbus/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dbus-cleanup-sockets.1.gz
usr/share/man/man1/dbus-daemon.1.gz
usr/share/man/man1/dbus-launch.1.gz
usr/share/man/man1/dbus-monitor.1.gz
usr/share/man/man1/dbus-run-session.1.gz
usr/share/man/man1/dbus-send.1.gz
usr/share/man/man1/dbus-test-tool.1.gz
usr/share/man/man1/dbus-update-activation-environment.1.gz
usr/share/man/man1/dbus-uuidgen.1.gz
usr/share/xml/
usr/share/xml/dbus-1/
usr/share/xml/dbus-1/busconfig.dtd
usr/share/xml/dbus-1/introspect.dtd
var/
var/lib/
var/lib/dbus/

