# alpm/tests/db/local/mkinitcpio-25-2/files
%FILES%
etc/
etc/initcpio/
etc/initcpio/hooks/
etc/initcpio/install/
etc/mkinitcpio.conf
etc/mkinitcpio.d/
usr/
usr/bin/
usr/bin/lsinitcpio
usr/bin/mkinitcpio
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/functions
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/consolefont
usr/lib/initcpio/hooks/keymap
usr/lib/initcpio/hooks/memdisk
usr/lib/initcpio/hooks/resume
usr/lib/initcpio/hooks/shutdown
usr/lib/initcpio/hooks/sleep
usr/lib/initcpio/hooks/usr
usr/lib/initcpio/init
usr/lib/initcpio/init_functions
usr/lib/initcpio/install/
usr/lib/initcpio/install/autodetect
usr/lib/initcpio/install/base
usr/lib/initcpio/install/block
usr/lib/initcpio/install/consolefont
usr/lib/initcpio/install/filesystems
usr/lib/initcpio/install/fsck
usr/lib/initcpio/install/hostdata
usr/lib/initcpio/install/keyboard
usr/lib/initcpio/install/keymap
usr/lib/initcpio/install/memdisk
usr/lib/initcpio/install/modconf
usr/lib/initcpio/install/resume
usr/lib/initcpio/install/sd-shutdown
usr/lib/initcpio/install/sd-vconsole
usr/lib/initcpio/install/shutdown
usr/lib/initcpio/install/sleep
usr/lib/initcpio/install/strip
usr/lib/initcpio/install/usr
usr/lib/initcpio/shutdown
usr/lib/initcpio/udev/
usr/lib/initcpio/udev/01-memdisk.rules
usr/lib/kernel/
usr/lib/kernel/install.d/
usr/lib/kernel/install.d/50-mkinitcpio.install
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mkinitcpio-generate-shutdown-ramfs.service
usr/lib/systemd/system/shutdown.target.wants/
usr/lib/systemd/system/shutdown.target.wants/mkinitcpio-generate-shutdown-ramfs.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mkinitcpio.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/lsinitcpio
usr/share/bash-completion/completions/mkinitcpio
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/lsinitcpio.1.gz
usr/share/man/man5/
usr/share/man/man5/mkinitcpio.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/mkinitcpio.8.gz
usr/share/mkinitcpio/
usr/share/mkinitcpio/example.preset
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_mkinitcpio

%BACKUP%
etc/mkinitcpio.conf	(null)

