# alpm/tests/db/local/iptables-1:1.8.2-1/files
%FILES%
etc/
etc/ethertypes
etc/iptables/
etc/iptables/empty.rules
etc/iptables/ip6tables.rules
etc/iptables/iptables.rules
etc/iptables/simple_firewall.rules
usr/
usr/bin/
usr/bin/arptables-nft
usr/bin/arptables-nft-restore
usr/bin/arptables-nft-save
usr/bin/ebtables-nft
usr/bin/ebtables-nft-restore
usr/bin/ebtables-nft-save
usr/bin/ip6tables
usr/bin/ip6tables-legacy
usr/bin/ip6tables-legacy-restore
usr/bin/ip6tables-legacy-save
usr/bin/ip6tables-nft
usr/bin/ip6tables-nft-restore
usr/bin/ip6tables-nft-save
usr/bin/ip6tables-restore
usr/bin/ip6tables-restore-translate
usr/bin/ip6tables-save
usr/bin/ip6tables-translate
usr/bin/iptables
usr/bin/iptables-legacy
usr/bin/iptables-legacy-restore
usr/bin/iptables-legacy-save
usr/bin/iptables-nft
usr/bin/iptables-nft-restore
usr/bin/iptables-nft-save
usr/bin/iptables-restore
usr/bin/iptables-restore-translate
usr/bin/iptables-save
usr/bin/iptables-translate
usr/bin/iptables-xml
usr/bin/nfbpf_compile
usr/bin/nfnl_osf
usr/bin/xtables-legacy-multi
usr/bin/xtables-monitor
usr/bin/xtables-nft-multi
usr/include/
usr/include/libipq.h
usr/include/libiptc/
usr/include/libiptc/ipt_kernel_headers.h
usr/include/libiptc/libip6tc.h
usr/include/libiptc/libiptc.h
usr/include/libiptc/libxtc.h
usr/include/libiptc/xtcshared.h
usr/include/xtables-version.h
usr/include/xtables.h
usr/lib/
usr/lib/libip4tc.so
usr/lib/libip4tc.so.0
usr/lib/libip4tc.so.0.1.0
usr/lib/libip6tc.so
usr/lib/libip6tc.so.0
usr/lib/libip6tc.so.0.1.0
usr/lib/libipq.so
usr/lib/libipq.so.0
usr/lib/libipq.so.0.0.0
usr/lib/libiptc.so
usr/lib/libiptc.so.0
usr/lib/libiptc.so.0.0.0
usr/lib/libxtables.so
usr/lib/libxtables.so.12
usr/lib/libxtables.so.12.2.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libip4tc.pc
usr/lib/pkgconfig/libip6tc.pc
usr/lib/pkgconfig/libipq.pc
usr/lib/pkgconfig/libiptc.pc
usr/lib/pkgconfig/xtables.pc
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/iptables-flush
usr/lib/systemd/system/
usr/lib/systemd/system/ip6tables.service
usr/lib/systemd/system/iptables.service
usr/lib/xtables/
usr/lib/xtables/libarpt_mangle.so
usr/lib/xtables/libebt_802_3.so
usr/lib/xtables/libebt_arp.so
usr/lib/xtables/libebt_arpreply.so
usr/lib/xtables/libebt_dnat.so
usr/lib/xtables/libebt_ip.so
usr/lib/xtables/libebt_ip6.so
usr/lib/xtables/libebt_log.so
usr/lib/xtables/libebt_mark.so
usr/lib/xtables/libebt_mark_m.so
usr/lib/xtables/libebt_nflog.so
usr/lib/xtables/libebt_pkttype.so
usr/lib/xtables/libebt_redirect.so
usr/lib/xtables/libebt_snat.so
usr/lib/xtables/libebt_stp.so
usr/lib/xtables/libebt_vlan.so
usr/lib/xtables/libip6t_DNAT.so
usr/lib/xtables/libip6t_DNPT.so
usr/lib/xtables/libip6t_HL.so
usr/lib/xtables/libip6t_LOG.so
usr/lib/xtables/libip6t_MASQUERADE.so
usr/lib/xtables/libip6t_NETMAP.so
usr/lib/xtables/libip6t_REDIRECT.so
usr/lib/xtables/libip6t_REJECT.so
usr/lib/xtables/libip6t_SNAT.so
usr/lib/xtables/libip6t_SNPT.so
usr/lib/xtables/libip6t_ah.so
usr/lib/xtables/libip6t_dst.so
usr/lib/xtables/libip6t_eui64.so
usr/lib/xtables/libip6t_frag.so
usr/lib/xtables/libip6t_hbh.so
usr/lib/xtables/libip6t_hl.so
usr/lib/xtables/libip6t_icmp6.so
usr/lib/xtables/libip6t_ipv6header.so
usr/lib/xtables/libip6t_mh.so
usr/lib/xtables/libip6t_rt.so
usr/lib/xtables/libip6t_srh.so
usr/lib/xtables/libipt_CLUSTERIP.so
usr/lib/xtables/libipt_DNAT.so
usr/lib/xtables/libipt_ECN.so
usr/lib/xtables/libipt_LOG.so
usr/lib/xtables/libipt_MASQUERADE.so
usr/lib/xtables/libipt_NETMAP.so
usr/lib/xtables/libipt_REDIRECT.so
usr/lib/xtables/libipt_REJECT.so
usr/lib/xtables/libipt_SNAT.so
usr/lib/xtables/libipt_TTL.so
usr/lib/xtables/libipt_ULOG.so
usr/lib/xtables/libipt_ah.so
usr/lib/xtables/libipt_icmp.so
usr/lib/xtables/libipt_realm.so
usr/lib/xtables/libipt_ttl.so
usr/lib/xtables/libxt_AUDIT.so
usr/lib/xtables/libxt_CHECKSUM.so
usr/lib/xtables/libxt_CLASSIFY.so
usr/lib/xtables/libxt_CONNMARK.so
usr/lib/xtables/libxt_CONNSECMARK.so
usr/lib/xtables/libxt_CT.so
usr/lib/xtables/libxt_DSCP.so
usr/lib/xtables/libxt_HMARK.so
usr/lib/xtables/libxt_IDLETIMER.so
usr/lib/xtables/libxt_LED.so
usr/lib/xtables/libxt_MARK.so
usr/lib/xtables/libxt_NFLOG.so
usr/lib/xtables/libxt_NFQUEUE.so
usr/lib/xtables/libxt_NOTRACK.so
usr/lib/xtables/libxt_RATEEST.so
usr/lib/xtables/libxt_SECMARK.so
usr/lib/xtables/libxt_SET.so
usr/lib/xtables/libxt_SYNPROXY.so
usr/lib/xtables/libxt_TCPMSS.so
usr/lib/xtables/libxt_TCPOPTSTRIP.so
usr/lib/xtables/libxt_TEE.so
usr/lib/xtables/libxt_TOS.so
usr/lib/xtables/libxt_TPROXY.so
usr/lib/xtables/libxt_TRACE.so
usr/lib/xtables/libxt_addrtype.so
usr/lib/xtables/libxt_bpf.so
usr/lib/xtables/libxt_cgroup.so
usr/lib/xtables/libxt_cluster.so
usr/lib/xtables/libxt_comment.so
usr/lib/xtables/libxt_connbytes.so
usr/lib/xtables/libxt_connlabel.so
usr/lib/xtables/libxt_connlimit.so
usr/lib/xtables/libxt_connmark.so
usr/lib/xtables/libxt_conntrack.so
usr/lib/xtables/libxt_cpu.so
usr/lib/xtables/libxt_dccp.so
usr/lib/xtables/libxt_devgroup.so
usr/lib/xtables/libxt_dscp.so
usr/lib/xtables/libxt_ecn.so
usr/lib/xtables/libxt_esp.so
usr/lib/xtables/libxt_hashlimit.so
usr/lib/xtables/libxt_helper.so
usr/lib/xtables/libxt_ipcomp.so
usr/lib/xtables/libxt_iprange.so
usr/lib/xtables/libxt_ipvs.so
usr/lib/xtables/libxt_length.so
usr/lib/xtables/libxt_limit.so
usr/lib/xtables/libxt_mac.so
usr/lib/xtables/libxt_mark.so
usr/lib/xtables/libxt_multiport.so
usr/lib/xtables/libxt_nfacct.so
usr/lib/xtables/libxt_osf.so
usr/lib/xtables/libxt_owner.so
usr/lib/xtables/libxt_physdev.so
usr/lib/xtables/libxt_pkttype.so
usr/lib/xtables/libxt_policy.so
usr/lib/xtables/libxt_quota.so
usr/lib/xtables/libxt_rateest.so
usr/lib/xtables/libxt_recent.so
usr/lib/xtables/libxt_rpfilter.so
usr/lib/xtables/libxt_sctp.so
usr/lib/xtables/libxt_set.so
usr/lib/xtables/libxt_socket.so
usr/lib/xtables/libxt_standard.so
usr/lib/xtables/libxt_state.so
usr/lib/xtables/libxt_statistic.so
usr/lib/xtables/libxt_string.so
usr/lib/xtables/libxt_tcp.so
usr/lib/xtables/libxt_tcpmss.so
usr/lib/xtables/libxt_time.so
usr/lib/xtables/libxt_tos.so
usr/lib/xtables/libxt_u32.so
usr/lib/xtables/libxt_udp.so
usr/share/
usr/share/iptables/
usr/share/iptables/empty-filter.rules
usr/share/iptables/empty-mangle.rules
usr/share/iptables/empty-nat.rules
usr/share/iptables/empty-raw.rules
usr/share/iptables/empty-security.rules
usr/share/iptables/empty.rules
usr/share/iptables/simple_firewall.rules
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/iptables-xml.1.gz
usr/share/man/man3/
usr/share/man/man3/ipq_create_handle.3.gz
usr/share/man/man3/ipq_destroy_handle.3.gz
usr/share/man/man3/ipq_errstr.3.gz
usr/share/man/man3/ipq_get_msgerr.3.gz
usr/share/man/man3/ipq_get_packet.3.gz
usr/share/man/man3/ipq_message_type.3.gz
usr/share/man/man3/ipq_perror.3.gz
usr/share/man/man3/ipq_read.3.gz
usr/share/man/man3/ipq_set_mode.3.gz
usr/share/man/man3/ipq_set_verdict.3.gz
usr/share/man/man3/libipq.3.gz
usr/share/man/man8/
usr/share/man/man8/ip6tables-restore.8.gz
usr/share/man/man8/ip6tables-save.8.gz
usr/share/man/man8/ip6tables.8.gz
usr/share/man/man8/iptables-extensions.8.gz
usr/share/man/man8/iptables-restore.8.gz
usr/share/man/man8/iptables-save.8.gz
usr/share/man/man8/iptables.8.gz
usr/share/man/man8/nfnl_osf.8.gz
usr/share/man/man8/xtables-legacy.8.gz
usr/share/man/man8/xtables-monitor.8.gz
usr/share/man/man8/xtables-nft.8.gz
usr/share/man/man8/xtables-translate.8.gz
usr/share/xtables/
usr/share/xtables/pf.os

%BACKUP%
etc/ethertypes	(null)
etc/iptables/iptables.rules	(null)
etc/iptables/ip6tables.rules	(null)

