# alpm/tests/db/local/krb5-1.16.1-1/files
%FILES%
etc/
etc/krb5.conf
usr/
usr/bin/
usr/bin/gss-client
usr/bin/gss-server
usr/bin/k5srvutil
usr/bin/kadmin
usr/bin/kadmin.local
usr/bin/kadmind
usr/bin/kdb5_ldap_util
usr/bin/kdb5_util
usr/bin/kdestroy
usr/bin/kinit
usr/bin/klist
usr/bin/kpasswd
usr/bin/kprop
usr/bin/kpropd
usr/bin/kproplog
usr/bin/krb5-config
usr/bin/krb5-send-pr
usr/bin/krb5kdc
usr/bin/ksu
usr/bin/kswitch
usr/bin/ktutil
usr/bin/kvno
usr/bin/sclient
usr/bin/sim_client
usr/bin/sim_server
usr/bin/sserver
usr/bin/uuclient
usr/bin/uuserver
usr/include/
usr/include/gssapi.h
usr/include/gssapi/
usr/include/gssapi/gssapi.h
usr/include/gssapi/gssapi_ext.h
usr/include/gssapi/gssapi_generic.h
usr/include/gssapi/gssapi_krb5.h
usr/include/gssapi/mechglue.h
usr/include/gssrpc/
usr/include/gssrpc/auth.h
usr/include/gssrpc/auth_gss.h
usr/include/gssrpc/auth_gssapi.h
usr/include/gssrpc/auth_unix.h
usr/include/gssrpc/clnt.h
usr/include/gssrpc/netdb.h
usr/include/gssrpc/pmap_clnt.h
usr/include/gssrpc/pmap_prot.h
usr/include/gssrpc/pmap_rmt.h
usr/include/gssrpc/rename.h
usr/include/gssrpc/rpc.h
usr/include/gssrpc/rpc_msg.h
usr/include/gssrpc/svc.h
usr/include/gssrpc/svc_auth.h
usr/include/gssrpc/types.h
usr/include/gssrpc/xdr.h
usr/include/kadm5/
usr/include/kadm5/admin.h
usr/include/kadm5/chpass_util_strings.h
usr/include/kadm5/kadm_err.h
usr/include/kdb.h
usr/include/krad.h
usr/include/krb5.h
usr/include/krb5/
usr/include/krb5/ccselect_plugin.h
usr/include/krb5/certauth_plugin.h
usr/include/krb5/clpreauth_plugin.h
usr/include/krb5/hostrealm_plugin.h
usr/include/krb5/kadm5_auth_plugin.h
usr/include/krb5/kadm5_hook_plugin.h
usr/include/krb5/kdcpolicy_plugin.h
usr/include/krb5/kdcpreauth_plugin.h
usr/include/krb5/krb5.h
usr/include/krb5/localauth_plugin.h
usr/include/krb5/locate_plugin.h
usr/include/krb5/plugin.h
usr/include/krb5/preauth_plugin.h
usr/include/krb5/pwqual_plugin.h
usr/include/profile.h
usr/include/verto-module.h
usr/include/verto.h
usr/lib/
usr/lib/krb5/
usr/lib/krb5/plugins/
usr/lib/krb5/plugins/kdb/
usr/lib/krb5/plugins/kdb/db2.so
usr/lib/krb5/plugins/kdb/kldap.so
usr/lib/krb5/plugins/preauth/
usr/lib/krb5/plugins/preauth/otp.so
usr/lib/krb5/plugins/preauth/pkinit.so
usr/lib/krb5/plugins/preauth/test.so
usr/lib/krb5/plugins/tls/
usr/lib/krb5/plugins/tls/k5tls.so
usr/lib/libgssapi_krb5.so
usr/lib/libgssapi_krb5.so.2
usr/lib/libgssapi_krb5.so.2.2
usr/lib/libgssrpc.so
usr/lib/libgssrpc.so.4
usr/lib/libgssrpc.so.4.2
usr/lib/libk5crypto.so
usr/lib/libk5crypto.so.3
usr/lib/libk5crypto.so.3.1
usr/lib/libkadm5clnt.so
usr/lib/libkadm5clnt_mit.so
usr/lib/libkadm5clnt_mit.so.11
usr/lib/libkadm5clnt_mit.so.11.0
usr/lib/libkadm5srv.so
usr/lib/libkadm5srv_mit.so
usr/lib/libkadm5srv_mit.so.11
usr/lib/libkadm5srv_mit.so.11.0
usr/lib/libkdb5.so
usr/lib/libkdb5.so.9
usr/lib/libkdb5.so.9.0
usr/lib/libkdb_ldap.so
usr/lib/libkdb_ldap.so.1
usr/lib/libkdb_ldap.so.1.0
usr/lib/libkrad.so
usr/lib/libkrad.so.0
usr/lib/libkrad.so.0.0
usr/lib/libkrb5.so
usr/lib/libkrb5.so.3
usr/lib/libkrb5.so.3.3
usr/lib/libkrb5support.so
usr/lib/libkrb5support.so.0
usr/lib/libkrb5support.so.0.1
usr/lib/libverto.so
usr/lib/libverto.so.0
usr/lib/libverto.so.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/gssrpc.pc
usr/lib/pkgconfig/kadm-client.pc
usr/lib/pkgconfig/kadm-server.pc
usr/lib/pkgconfig/kdb.pc
usr/lib/pkgconfig/krb5-gssapi.pc
usr/lib/pkgconfig/krb5.pc
usr/lib/pkgconfig/mit-krb5-gssapi.pc
usr/lib/pkgconfig/mit-krb5.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/krb5-kadmind.service
usr/lib/systemd/system/krb5-kdc.service
usr/lib/systemd/system/krb5-kpropd.service
usr/lib/systemd/system/krb5-kpropd.socket
usr/lib/systemd/system/krb5-kpropd@.service
usr/share/
usr/share/aclocal/
usr/share/aclocal/ac_check_krb5.m4
usr/share/doc/
usr/share/doc/krb5/
usr/share/doc/krb5/examples/
usr/share/doc/krb5/examples/kdc.conf
usr/share/doc/krb5/examples/kerberos.ldif
usr/share/doc/krb5/examples/kerberos.schema
usr/share/doc/krb5/examples/krb5.conf
usr/share/doc/krb5/examples/services.append
usr/share/licenses/
usr/share/licenses/krb5/
usr/share/licenses/krb5/LICENSE
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/mit-krb5.mo
usr/share/locale/en_US/
usr/share/locale/en_US/LC_MESSAGES/
usr/share/locale/en_US/LC_MESSAGES/mit-krb5.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/k5srvutil.1.gz
usr/share/man/man1/kadmin.1.gz
usr/share/man/man1/kdestroy.1.gz
usr/share/man/man1/kinit.1.gz
usr/share/man/man1/klist.1.gz
usr/share/man/man1/kpasswd.1.gz
usr/share/man/man1/krb5-config.1.gz
usr/share/man/man1/ksu.1.gz
usr/share/man/man1/kswitch.1.gz
usr/share/man/man1/ktutil.1.gz
usr/share/man/man1/kvno.1.gz
usr/share/man/man1/sclient.1.gz
usr/share/man/man5/
usr/share/man/man5/.k5identity.5.gz
usr/share/man/man5/.k5login.5.gz
usr/share/man/man5/k5identity.5.gz
usr/share/man/man5/k5login.5.gz
usr/share/man/man5/kadm5.acl.5.gz
usr/share/man/man5/kdc.conf.5.gz
usr/share/man/man5/krb5.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/kadmin.local.8.gz
usr/share/man/man8/kadmind.8.gz
usr/share/man/man8/kdb5_ldap_util.8.gz
usr/share/man/man8/kdb5_util.8.gz
usr/share/man/man8/kprop.8.gz
usr/share/man/man8/kpropd.8.gz
usr/share/man/man8/kproplog.8.gz
usr/share/man/man8/krb5kdc.8.gz
usr/share/man/man8/sserver.8.gz
var/
var/lib/
var/lib/krb5kdc/
var/lib/krb5kdc/kdc.conf

%BACKUP%
etc/krb5.conf	(null)
var/lib/krb5kdc/kdc.conf	(null)

