# alpm/tests/db/local/device-mapper-2.02.183-2/files
%FILES%
usr/
usr/bin/
usr/bin/blkdeactivate
usr/bin/dmeventd
usr/bin/dmsetup
usr/bin/dmstats
usr/include/
usr/include/libdevmapper-event.h
usr/include/libdevmapper.h
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/udev/
usr/lib/initcpio/udev/11-dm-initramfs.rules
usr/lib/libdevmapper-event.so
usr/lib/libdevmapper-event.so.1.02
usr/lib/libdevmapper.so
usr/lib/libdevmapper.so.1.02
usr/lib/pkgconfig/
usr/lib/pkgconfig/devmapper-event.pc
usr/lib/pkgconfig/devmapper.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dm-event.service
usr/lib/systemd/system/dm-event.socket
usr/lib/systemd/system/sockets.target.wants/
usr/lib/systemd/system/sockets.target.wants/dm-event.socket
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/10-dm.rules
usr/lib/udev/rules.d/13-dm-disk.rules
usr/lib/udev/rules.d/95-dm-notify.rules
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/blkdeactivate.8.gz
usr/share/man/man8/dmeventd.8.gz
usr/share/man/man8/dmsetup.8.gz
usr/share/man/man8/dmstats.8.gz

